# gws/app/supervisor_t.conf
[unix_http_server]
file=/tmp/supervisor.sock                                 ; path to your socket file
chmod = 0777
username = user
password = 123

[supervisord]
pidfile=/var/supervisord/supervisord.pid                ; pidfile location
nodaemon=false                                            ; run supervisord as a daemon
minfds=1024                                               ; number of startup file descriptors
minprocs=200                                              ; number of process descriptors
user=USER ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock                     ; use a unix:// URL  for a unix socket

[program:ssh-aws]
process_name=%(program_name)s_%(process_num)02d
command=sh start_socks.sh 128%(process_num)02d
directory = START_SOCKS_DIR
autostart=true
autorestart=true
user=USER
numprocs=10
redirect_stderr=true

[program:ssh-nginx]
command=NGINX_BASE_DIR/sbin/nginx
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
